# baked-no-volume/supervisor/supervisord.conf
; supervisord — starts and keeps alive every service in the image.
; Logs go to /var/log/<service>.log (view with `cat` or in JupyterLab).

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.log

[program:comfyui]
command=bash /workspace/scripts/services/comfyui.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/comfyui.log
stderr_logfile=/var/log/comfyui.log

[program:kohya]
command=bash /workspace/scripts/services/kohya.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/kohya.log
stderr_logfile=/var/log/kohya.log

[program:filebrowser]
command=bash /workspace/scripts/services/filebrowser.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/filebrowser.log
stderr_logfile=/var/log/filebrowser.log

[program:jupyter]
command=bash /workspace/scripts/services/jupyter.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/jupyter.log
stderr_logfile=/var/log/jupyter.log

[program:codeserver]
command=bash /workspace/scripts/services/code-server.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/codeserver.log
stderr_logfile=/var/log/codeserver.log
